# docker/s6-rc.d/init-app/run
#!/command/with-contenv sh
set -eu

# Wait for Postgres to accept connections
until pg_isready -h 127.0.0.1 -U postgres >/dev/null 2>&1; do
  sleep 0.5
done

# Create the application database if it doesn't exist
psql -h 127.0.0.1 -U postgres -tc "SELECT 1 FROM pg_database WHERE datname='logbook'" \
  | grep -q 1 \
  || psql -h 127.0.0.1 -U postgres -c "CREATE DATABASE logbook"

# Run migrations (embedded SQL files in /app/migrations)
cd /app
for f in migrations/*.sql; do
  [ -f "$f" ] && psql -h 127.0.0.1 -U postgres -d logbook -v ON_ERROR_STOP=1 -f "$f"
done
