# docker/s6-rc.d/postgres/run
#!/command/with-contenv sh
exec su-exec postgres postgres -D "$PGDATA"
